# supervisord.conf
[supervisord]
nodaemon=true

[program:socket]
command=python /app/src/socket_runner.py

[program:graph]
command=python /app/src/graph_runner.py